# systemd/.config/systemd/user/gnome-keyring.service
[Unit]
Description=GNOME Keyring Daemon
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/gnome-keyring-daemon --foreground --components=secrets,pkcs11,ssh
Restart=on-failure
Slice=session.slice

[Install]
WantedBy=graphical-session.target
